# roles/cleanup/files/scheduled-reboot.service
[Unit]
Description=Scheduled Reboot

[Service]
Type=simple
ExecStart=/sbin/reboot
